# server/etc/cron.daily/article_rename_script
#!/bin/bash


/root/bin/article_renamer.py -d /mnt/data3/movie
/root/bin/article_renamer.py -d /mnt/data2/tv

