# server/etc/cron.daily/daily_backup
#!/bin/bash


#/mnt/backup
if mountpoint -q /mnt/data; then
    echo "Archiving /mnt/data"
    date >> /mnt/backup/daily/cronlog.daily
    rsync --archive --delete /mnt/data /mnt/backup/daily &>> /mnt/backup/daily/cronlog.daily
    rsync --archive --delete /etc /mnt/backup/daily &>> /mnt/backup/daily/cronlog.daily
else
    echo "/mnt/data is not mounted"
fi

#/mnt/backup2
if mountpoint -q /mnt/data2; then
    echo "Archiving /mnt/data2"
    date >> /mnt/backup2/daily/cronlog.daily
    rsync --archive --delete /mnt/data2 /mnt/backup2/daily &>> /mnt/backup2/daily/cronlog.daily
else
    echo "/mnt/data2 is not mounted"
fi

#/mnt/backup3
if mountpoint -q /mnt/data3; then
    echo "Archiving /mnt/data3"
    date >> /mnt/backup3/daily/cronlog.daily
    rsync --archive --delete /mnt/data3 /mnt/backup3/daily &>> /mnt/backup3/daily/cronlog.daily
else
    echo "/mnt/data3 is not mounted"
fi
